# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/serverdata/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:ark

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1024x768x16
user=ark
autostart=true
autorestart=true
priority=10
stdout_logfile=/serverdata/logs/xvfb.log
stderr_logfile=/serverdata/logs/xvfb.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:arkserver]
command=/app/scripts/start_ark.sh
user=ark
autostart=false
autorestart=false
priority=20
stopwaitsecs=60
environment=DISPLAY=":1",STEAM_COMPAT_DATA_PATH="/serverdata/proton-compat",STEAM_COMPAT_CLIENT_INSTALL_PATH="/opt/steamcmd",HOME="/home/ark",PROTON_USE_WINED3D="1",PROTON_LOG="1",WINEDLLOVERRIDES="version=n,b",SteamAppId="2430930",SteamGameId="2430930"
stdout_logfile=/serverdata/logs/arkserver.log
stderr_logfile=/serverdata/logs/arkserver.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

[program:webui]
command=/app/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --timeout-graceful-shutdown 85
user=ark
directory=/app
autostart=true
autorestart=true
priority=30
stopwaitsecs=90
stdout_logfile=/serverdata/logs/webui.log
stderr_logfile=/serverdata/logs/webui.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
